# root/etc/services.d/natpmp/run
#!/usr/bin/with-contenv bash

/usr/bin/transmission-rust
