# root/etc/services.d/natpmp/run
#!/usr/bin/with-contenv bash

/user/bin/natpmp-setup